# deploy/systemd/ttc-workclock-api.service
[Unit]
Description=TTC WorkClock API
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=ttc-workclock
Group=ttc-workclock
WorkingDirectory=/opt/ttc-workclock/backend
EnvironmentFile=/opt/ttc-workclock/.env
Environment=PYTHONPATH=/opt/ttc-workclock/backend
ExecStart=/opt/ttc-workclock/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
